# templates/etc/systemd/td-agent.service.erb
[Unit]
Description=<%= project_name %>: Fluentd based data collector for Treasure Data
Documentation=https://docs.treasuredata.com/articles/td-agent
After=network-online.target
Wants=network-online.target

[Service]
UMask=027
User=<%= Shellwords.shellescape(project_name) %>
Group=<%= Shellwords.shellescape(project_name) %>
LimitNOFILE=65536
Environment=LD_PRELOAD=<%= install_path %>/embedded/lib/libjemalloc.so
Environment=GEM_HOME=<%= install_path %>/embedded/lib/ruby/gems/<%= gem_dir_version %>/
Environment=GEM_PATH=<%= install_path %>/embedded/lib/ruby/gems/<%= gem_dir_version %>/
Environment=FLUENT_CONF=/etc/<%= project_name %>/<%= project_name %>.conf
Environment=FLUENT_PLUGIN=/etc/<%= project_name %>/plugin
Environment=FLUENT_SOCKET=/var/run/<%= project_name %>/<%= project_name %>.sock
Environment=TD_AGENT_LOG_FILE=/var/log/<%= project_name %>/<%= project_name %>.log
Environment=TD_AGENT_OPTIONS=
<% if pkg_type == 'deb' %>
EnvironmentFile=-/etc/default/<%= project_name %>
<% else %>
EnvironmentFile=-/etc/sysconfig/<%= project_name %>
<% end %>
PIDFile=<%= Shellwords.shellescape(File.join(root_path, "var", "run", project_name, "#{project_name}.pid")) %>
RuntimeDirectory=<%= Shellwords.shellescape(project_name) %>
Type=forking
ExecStart=/opt/td-agent/embedded/bin/fluentd --log $TD_AGENT_LOG_FILE --daemon <%= Shellwords.shellescape(File.join(root_path, "var", "run", project_name, "#{project_name}.pid")) %> $TD_AGENT_OPTIONS
ExecStop=/bin/kill -TERM ${MAINPID}
ExecReload=/bin/kill -HUP ${MAINPID}
Restart=always
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
