# templates/package-scripts/td-agent/deb/prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/<%= project_name %>" ]; then
	invoke-rc.d <%= project_name %> stop || exit $?
fi

