# templates/td-agent.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>KeepAlive</key>
  <true/>
  <key>Label</key>
  <string><%= project_name %></string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/sbin/<%= project_name %></string>
    <string>--log</string>
    <string>/var/log/<%= project_name %>/<%= project_name %>.log</string>
    <string>--use-v1-config</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>WorkingDirectory</key>
  <string><%= install_path %></string>
</dict>
</plist>
